# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/sbin/zebra -d -u root -g root -f /etc/quagga/zebra.conf
[program:bgpd]
command=/usr/sbin/bgpd -d -u root -g root -f /etc/quagga/bgpd.conf
